# systemd/gost.service
[Unit]
Description=gost tunnel service
After=network.target

[Service]
Type=simple
User=gostsvc
Group=gostsvc
ExecStart=/usr/local/bin/gost -C /etc/gost/config.json
Restart=on-failure
LimitNOFILE=65536

[Install]
WantedBy=multi-user.target
